# dist/debian/arm64/usr/lib/systemd/user/smolsonic.service
[Unit]
Description=smolsonic Subsonic-compatible music server
After=network.target

[Service]
ExecStart=/usr/local/bin/smolsonic --config %h/.config/smolsonic/smolsonic.toml
WorkingDirectory=%h/.local/share/smolsonic
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
